# deploy/partner-edge/systemd/oxpulse-partner-edge-hydrate.service
[Unit]
Description=OxPulse partner-edge per-clone hydration
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service
Before=oxpulse-partner-edge.service
ConditionPathExists=!/var/lib/oxpulse-partner-edge/hydrated

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/sbin/oxpulse-partner-edge-hydrate
StandardOutput=journal
StandardError=journal
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target
